# deploy/cron/ingest.crontab
# Sync legacy drop from deathstar, then run ingest every 5 minutes
*/5 * * * * cd /app && /usr/bin/env bash deploy/sync_from_legacy_drop.sh >> /var/log/cron-sync.log 2>&1 && /usr/bin/env bash run_ingest_all.sh >> /var/log/cron-ingest.log 2>&1
